# gcp-image-from-docker/builtin/dhclient.service
[Unit]
Description=DHCP Client
Documentation=man:dhclient(8)
Wants=network.target
Before=network.target

[Service]
Type=simple
PIDFile=/var/run/dhclient.pid
ExecStart=/sbin/dhclient -d

[Install]
WantedBy=multi-user.target
